# system/etc/runit/sv/ntpd/run
#!/bin/sh
# Thanks to replabrobin for the network check!
[ "$(ip route | awk '/^default/{print ($3!="")+($5!="")}')" = "2" ] || exit 1
# wait for a usable global address, not just a default route (ntpd segfaults during boot-time address churn)
ip route get 1.1.1.1 >/dev/null 2>&1 || exit 1
exec ntpd -g -u ntp:ntp -n >/dev/null 2>&1
